# weewx-ha/rootfs/etc/s6-overlay/s6-rc.d/.weewx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: SQLite Web
# Runs SQLite Web
# ==============================================================================

bashio::log.info 'WeeWX-HA...'

exec /home/weewx/bin/weewxd /config/weewx/weewx.conf
