# weewx5-full/rootfs/etc/s6-overlay/s6-rc.d/gunicorn/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: SQLite Web
# Take down the S6 supervision tree when Gunicorn fails
# ==============================================================================
