# weewx5-full/rootfs/etc/s6-overlay/s6-rc.d/gunicorn/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: WeeWX-HA
# Runs the Gunicorn daemon
# ==============================================================================

bashio::log.info "Starting gunicorn..."

# Start the python web server
source ~/gunicorn-venv/bin/activate
#source /root/gunicorn-venv/bin/activate

cd /var/www/python
# gunicorn  --bind unix:/run/wsgi.sock --bind 0.0.0.0:5000 'wsgi:app'
gunicorn  --bind unix:/run/wsgi.sock 'wsgi:app'