# weewx5/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: WeeWX-HA
# Configures NGINX for use with WeeWX
# ==============================================================================
set -e
CONFIG_PATH=/data/options.json
WEEWX_CFG=/config/weewx.conf

# Create the html directory if needed
cp -n -r /var/www/html /config

# Generate Ingress configuration
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf
    
# Generate direct access configuration, if enabled.
if bashio::var.has_value "$(bashio::addon.port 80)"; then
    tempio \
        -template /etc/nginx/templates/direct.gtpl \
        -out /etc/nginx/servers/direct.conf
fi  